# pibags-1.0.0/DEBIAN/postinst
#!/bin/sh
pip install osc4py3
cp /home/pi/pibags/pibags.service /etc/systemd/system/
systemctl daemon-reload
systemctl enable pibags
